# imepp-agenda.service
[Unit]
Descriptio=IME++ Agenda

[Service]
ExecStart=/srv/imepp-agenda/run.sh

[Install]
WantedBy=multi-user.target
